# corona.service
[Service]
User=corona
KillMode=process
ExecStart=/bin/bash -lc 'exec /home/corona/current/run'

[Install]
WantedBy=multi-user.target
