# corona.service
[Service]
User=corona
ExecStart=/bin/bash -lc '/home/corona/current/run'

[Install]
WantedBy=multi-user.target
